# smf.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='local/auditrecv' type='service' version='0'>
    <create_default_instance enabled='true'/>
    <single_instance/>
    <dependency restart_on="none" type="service"
            name="multi_user_dependency" grouping="require_all">
            <service_fmri value="svc:/milestone/multi-user"/>
        </dependency>
    <dependency name='fs-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='network' grouping='require_all' restart_on='restart' type='service'>
      <service_fmri value='svc:/milestone/network'/>
    </dependency>
    <method_context>
    <method_credential user='auditrecv' group='auditrecv' />
  <method_environment>
  <envvar name='PATH' value='/usr/local/sbin:/usr/local/bin:/opt/local/sbin:/opt/local/bin:/usr/sbin:/usr/bin:/sbin'/>
  <envvar name='HOME' value='/var/db/auditrecv'/>
  <envvar name='LOGNAME' value='auditrecv'/>
  <envvar name='RUNNER_LOG_DIR' value='/var/log/auditrecv'/>
  </method_environment>
  </method_context>
    <exec_method name='start' type='method' exec='/opt/auditrecv/bin/auditrecv foreground &amp;' timeout_seconds='60'/>
    <exec_method name='stop' type='method' exec='/opt/auditrecv/bin/auditrecv stop' timeout_seconds='60'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>audit receiver</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
